# container-files/etc/services.d/keepalived/run
#!/bin/execlineb -P

/usr/sbin/keepalived -ln
